# infra/ansible/templates/forge-agent.service.j2
[Unit]
Description=Forge worker agent
After=network-online.target
Wants=network-online.target

[Service]
User=forge
Group=forge
Environment=FORGE_CONTROL_PLANE_URL={{ forge_control_plane_url }}
Environment=FORGE_RUNNER_PATH=/usr/local/bin/forge-build-runner
Environment=FORGE_METRICS_SOCKET=/run/forge-agent/metrics.sock
Environment=FORGE_AGENT_APP_ROOT=/var/lib/forge-agent/apps
EnvironmentFile=-/etc/forge/agent.env
ExecStart=/usr/local/bin/forge-agent
Restart=always
RestartSec=3
UMask=0007
RuntimeDirectory=forge-agent
RuntimeDirectoryMode=0755
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/forge-agent /run/forge-agent /tmp
CapabilityBoundingSet=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
IPAddressDeny=169.254.169.254/32
IPAddressDeny=fd00:ec2::254/128
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
